# indicator/90-charging_enabled.rules
# indicator-upower

SUBSYSTEM=="power_supply", TEST=="/sys/class/power_supply/battery/charging_enabled", \
    RUN+="/bin/chmod 0664 /sys/class/power_supply/battery/charging_enabled"

SUBSYSTEM=="power_supply", TEST=="/sys/class/power_supply/battery/charging_enabled", \
    RUN+="/bin/chown phablet /sys/class/power_supply/battery/charging_enabled"


SUBSYSTEM=="power_supply", TEST=="/sys/class/power_supply/battery/battery_charging_enabled", \
    RUN+="/bin/chmod 0664 /sys/class/power_supply/battery/battery_charging_enabled"

SUBSYSTEM=="power_supply", TEST=="/sys/class/power_supply/battery/battery_charging_enabled", \
    RUN+="/bin/chown phablet /sys/class/power_supply/battery/battery_charging_enabled"
